# update-system-notify/cfg/update-system-notify.cron
#!/bin/sh

/usr/sbin/update-system-notify-daemon.sh &>/dev/null

